# debian/leela.postrm
#!/bin/sh

update-rc.d leela remove
update-rc.d leela-collectd remove
update-rc.d leela-dmproc remove
update-rc.d leela-http remove
update-rc.d leela-multicast remove
update-rc.d leela-storage remove
update-rc.d leela-timeline remove
update-rc.d leela-udp remove
update-rc.d leela-xmpp remove


exit 0